# grommunio-web/build-assets/s6-overlay/s6-rc.d/config-setup/up
#!/command/execlineb -P
# "a oneshot script must be execline because of how s6-rc internally works"
/etc/s6-overlay/s6-rc.d/config-setup/config-setup.sh
