# build-assets/services/config-setup/run
#!/command/execlineb -P

mv /etc/php7/fpm/php-fpm.conf.default /etc/php7/fpm/php-fpm.conf  
ln -s /home/gromox-services/* /etc/gromox/ 
