# build-assets/services/gromox-antispam/run
#!/command/execlineb -P

mkdir -p /var/run/grommunio-antispam 
rspamd -u groas -g grommunio 
redis-server
