# grommunio-core/build-assets/s6-overlay/s6-rc.d/php-fpm/run
#!/bin/sh -e
exec 2>&1
if [ ! -S /run/php-fpm/php-gromox-fpm.sock ]; then
	mkdir -p /run/php-fpm/
	/usr/sbin/php-fpm --nodaemonize
fi
