# build-assets/s6-overlay/s6-rc.d/gromox-antispam-log/run
#!/bin/execlineb -P
s6-log -b p[antispam] T 1
