# build-assets/s6-overlay/s6-rc.d/gromox-antispam-log/run
#!/command/execlineb -P
s6-envuidgid -D 65534:65534 nobody
s6-applyuidgid -U
exec -c
s6-log -b -- T p[antispam] 1
