# build-assets/services/gromox-event/run
#!/command/execlineb -P

/usr/libexec/gromox/event
