# build-assets/s6-overlay/s6-rc.d/gromox-postfix/run
#!/command/execlineb -P
foreground { postfix check }
fdmove -c 2 1
exec /usr/lib/postfix/bin/master -d
