# build-assets/services/gromox-chat/run
#!/command/execlineb -P

cd /usr/share/grommunio-chat 
/home/chat/setup.sh 
