# build-assets/services/db-setup/run
#!/command/execlineb -P

/scripts/db.sh
