# build-assets/services/db-setup/up
#!/usr/bin/with-contenv sh
/scripts/db.sh
