# build-assets/services/gromox-delivery/run
#!/command/execlineb -P

/usr/libexec/gromox/delivery
