# build-assets/s6-overlay/s6-rc.d/gromox-admin/run
#!/command/execlineb -P
fdmove -c 2 1
uwsgi --uid grommunio --ini /usr/share/grommunio-admin-api/api-config.ini 
