# build-assets/s6-overlay/s6-rc.d/php-fpm/run
#!/bin/sh -e

if [ ! -S /run/gromox/php-fpm.sock ]; then
	/usr/sbin/php-fpm
fi
