# build-assets/s6-overlay/s6-rc.d/gromox-antispam/run
#!/command/execlineb -P

rspamd -u groas -g grommunio 
