# grommunio-core/build-assets/s6-overlay/s6-rc.d/gromox-event-log/run
#!/command/execlineb -P
/etc/s6-overlay/common/log-prefix p[event]
