# build-assets/services/gromox-pop3/run
#!/command/execlineb -P

/usr/libexec/gromox/pop3
