# build-assets/services/php-fpm/run
#!/command/execlineb -P

mv /etc/php7/fpm/php-fpm.conf.default /etc/php7/fpm/php-fpm.conf  
php-fpm 
