# build-assets/s6-overlay/s6-rc.d/gromox-http/run
#!/command/execlineb -P
exec 2>&1
/usr/libexec/gromox/http
